# templates/rclone-notify-error.target.j2
[Unit]
Description=Rclone error notifier
StopWhenUnneeded=true

[Install]
WantedBy=multi-user.target {{ rclone_service_name }}.target

